# php82-nginx-1.0/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm -F
priority=10
autostart=true
autorestart=true

[program:nginx]
command=nginx -g "daemon off;"
priority=20
autostart=true
autorestart=true
