# laravel11-docker/docker-config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
priority=5
startsecs=1
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
priority=10
startsecs=1
autorestart=true